# packaging/linux/tuxbot/nightly.service
[Unit]
Description=make a nightly build via keybase chat

[Service]
ExecStart=bash -c "keybase chat send keybase.bots.build.linux '!nightly'"
